# infra/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

; ── FastAPI ──────────────────────────────────────────────────────────────────
[program:fastapi]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2
directory=/app
user=appuser
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/supervisord/fastapi.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/supervisord/fastapi.log
environment=PYTHONPATH="/app"

; ── Celery Worker（多队列合并，4 并发）───────────────────────────────────────
[program:celery-worker]
command=celery -A app.workers.celery_app worker
    -Q agent_long,embedding,notify,vault_sync
    --concurrency=4
    --loglevel=info
    --logfile=/var/log/supervisord/celery.log
directory=/app
user=appuser
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/var/log/supervisord/celery.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/supervisord/celery.log
environment=PYTHONPATH="/app"
stopasgroup=true
killasgroup=true

; ── Celery Beat（定时任务调度，可选）─────────────────────────────────────────
[program:celery-beat]
command=celery -A app.workers.celery_app beat --loglevel=info
directory=/app
user=appuser
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisord/celery-beat.log
stderr_logfile=/var/log/supervisord/celery-beat.log
environment=PYTHONPATH="/app"

; ── Vault Watcher（长驻文件监听进程）────────────────────────────────────────
[program:vault-watcher]
command=python -m app.workers.vault_watcher_main
directory=/app
user=appuser
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/var/log/supervisord/vault-watcher.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisord/vault-watcher.log
environment=PYTHONPATH="/app"

; ── MCP SSE Server（可选，独立端口 8765）────────────────────────────────────
[program:mcp-server]
command=uvicorn app.agents.mcp_server:sse_app --host 0.0.0.0 --port 8765
directory=/app
user=appuser
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisord/mcp-server.log
stderr_logfile=/var/log/supervisord/mcp-server.log
environment=PYTHONPATH="/app"

; ── 进程组：一键管理所有必须进程 ──────────────────────────────────────────────
[group:lumipath]
programs=fastapi,celery-worker,vault-watcher
priority=999
